# launchd/com.fpof.market-intel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.fpof.market-intel</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/sherman/07. FPOF V2.2 Claude/scripts/market-intel/run-pipeline.sh</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>1</integer>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>WorkingDirectory</key>
  <string>/Users/sherman/07. FPOF V2.2 Claude</string>

  <key>StandardOutPath</key>
  <string>/Users/sherman/07. FPOF V2.2 Claude/data/market-intel/launchd-stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/sherman/07. FPOF V2.2 Claude/data/market-intel/launchd-stderr.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
  </dict>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
